# miniconda3/etc/profile.d/conda.csh
setenv CONDA_EXE "/data2/ai3625/workspace/ai3625-18/miniconda3/bin/conda";
setenv _CONDA_EXE "/data2/ai3625/workspace/ai3625-18/miniconda3/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/data2/ai3625/workspace/ai3625-18/miniconda3/bin/python";
setenv _CONDA_ROOT "/data2/ai3625/workspace/ai3625-18/miniconda3";

source "/data2/ai3625/workspace/ai3625-18/miniconda3/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"